# net/haproxy/files/etc/uci-defaults/haproxy
#!/bin/sh

rm -f /tmp/luci-*
[ ! -f /etc/init.d/haproxy ] || /etc/init.d/haproxy stop
[ ! -f /etc/init.d/haproxy ] || /etc/init.d/haproxy disable

exit 0
